# tmp/crontab.QapHYX/crontab
# access using `sudo crontab -e` 
# use the root user

* * * * * root /usr/local/bin/wifi-up.sh
