# tmp/crontab.QapHYX/crontab
* * * * * root /usr/local/bin/wifi-up.sh
